# systemd/waypaper-restore.service
[Unit]
Description=Restore Waypaper wallpapers on boot
PartOf=niri.service
After=niri.service dms.service
Requisite=niri.service

[Service]
Type=oneshot
# Wayland session environment - REQUIRED for linux-wallpaperengine to work
Environment=XDG_SESSION_TYPE=wayland
Environment=WAYLAND_DISPLAY=wayland-1
ExecStart=%h/.local/bin/waypaper-video-random --restore-or-random --restore-only
Restart=on-failure
RestartSec=5

[Install]
WantedBy=niri.service
